# images/nostalgia-crt/firstboot/nostalgia-firstboot.service
[Unit]
Description=Nostalgia OS first-boot setup
After=network-online.target
Wants=network-online.target
Before=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/lib/nostalgia/10-firstboot.sh
RemainAfterExit=yes
TimeoutStopFailureMode=abort
ConditionPathExists=!/var/lib/nostalgia/firstboot.done

[Install]
WantedBy=graphical.target
