# images/nostalgia-crt/firstboot/nostalgia-firstboot.service
[Unit]
Description=Nostalgia OS first-boot setup
Wants=network-online.target
After=network-online.target display-manager.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/sysconfig/nostalgia-firstboot
ExecStart=/usr/lib/nostalgia/10-nostalgia-setup.sh

[Install]
WantedBy=graphical.target
