# images/nostalgia-crt/firstboot/nostalgia.service
[Unit]
Description=Nostalgia OS first-boot setup
After=network-online.target
Wants=network-online.target
# Make sure we run our setup before the greeter so autologin applies immediately
Before=display-manager.service graphical.target

[Service]
Type=oneshot
ExecStart=/etc/firstboot.d/10-nostalgia-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
