# images/nostalgia-crt/firstboot/nostalgia.service
[Unit]
Description=First boot customization for Nostalgia CRT
After=network-online.target
Before=graphical.target

[Service]
Type=oneshot
ExecStart=/etc/firstboot.d/10-nostalgia-setup.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
