# images/nostalgia-arcade/systemd/nostalgia-setup.service
[Unit]
Description=Nostalgia Arcade First Boot Setup
After=default.target
ConditionUser=!@system

[Service]
Type=oneshot
ExecStart=/usr/libexec/nostalgia-setup
RemainAfterExit=yes

[Install]
WantedBy=default.target
